# deploy/systemd/cc.service.template
[Unit]
Description=Postline (CC) — Claude on Feishu, 24/7
Documentation=https://github.com/Christianye/postline
After=network-online.target
Wants=network-online.target
# Backoff across crashes.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{REPO_DIR}}

# Load env (credentials, allowlist, memory paths, etc.)
EnvironmentFile={{CC_HOME}}/env

# Absolute path to node — avoid relying on login-shell PATH.
ExecStart={{NODE_BIN}} packages/cli/dist/bin.js feishu

# Restart policy.
Restart=on-failure
RestartSec=5s

# Resource hygiene.
LimitNOFILE=4096

# Hardening — lightweight, avoids breaking things like /tmp usage.
PrivateTmp=false
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=false
# Logs via journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cc

[Install]
WantedBy=multi-user.target
